# discord-bot.ini
[program:discord-bot]
user = bbt
directory=/home/bbt/discord-bot
command=source /home/bbt/.virtualenvs/bin/postactivate && /home/bbt/.virtualenvs/discord/bin/python main.py
autostart=true
autorestart=true
redirect_stderr=true
stderr_logfile = /home/bbt/discord-bot/log/err.log
stdout_logfile = /home/bbt/discord-bot/log/out.log
